# supervisord.conf
[supervisord]
nodaemon=true

[program:chc]
command=/usr/sbin/chef-check-hostname

[program:chef-server]
command=/opt/chef-server/embedded/bin/runsvdir-start
